# ops/systemd/bostonmesh-dev.service
[Unit]
Description=Boston Mesh Docusaurus dev host
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yellowcooln
Group=yellowcooln
WorkingDirectory=/home/yellowcooln/bostonmesh
Environment=HOST=0.0.0.0
Environment=PORT=3000
Environment=BROWSER=none
ExecStart=/home/yellowcooln/bostonmesh/scripts/run-dev-host.sh
Restart=always
RestartSec=5
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
